# deploy/com.murakumo.comfyui.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!-- com.murakumo.comfyui — resident ComfyUI on a fleet COMPUTE node.

     LaunchDaemon (SYSTEM domain, /Library/LaunchDaemons + `UserName`), same
     shape as this repo's own com.murakumo.kotoba-mesh.plist.tmpl — NOT a
     per-user LaunchAgent. Two other approaches were tried and verified NOT
     to work on these nodes before landing on this one (2026-07-13, live on
     `benjamin`+7 other fleet nodes):

       1. gui/$(id -u) LaunchAgent — `launchctl bootstrap gui/<uid> ...`
          fails ("125: Domain does not support specified action") because
          these fleet Mac minis normally have NO active console session
          (`who` empty, `launchctl print gui/<uid>` itself unreachable).
       2. A `* * * * *` cron watchdog — `crontab <file>` (any form: plain,
          piped from `crontab -l`, even `sudo crontab -u <user> <file>`)
          fails with "Operation not permitted" over SSH: macOS gates
          *writing* a crontab behind Full Disk Access TCC consent, which an
          SSH-spawned shell doesn't have and can't grant itself non-
          interactively. (Reading — `crontab -l` — is unaffected, which is
          why the pre-existing `@reboot` line could be discovered but not
          replaced.) This is also why the OLD `@reboot ... comfyui/main.py`
          crontab line on each node could NOT be removed as part of this
          fix — it's still there, and will race this LaunchDaemon for
          :8188 at the next actual physical reboot. Harmless: whichever
          loses the bind just fails and (for the LaunchDaemon) retries every
          ThrottleInterval; only fix is removing that line at the physical
          console (or granting sshd/Terminal Full Disk Access), not scoped
          into this change.

     The system domain has no console-session dependency (`sudo launchctl
     bootstrap system ...` works headless) — verified live: killed the
     running process with `sudo kill -9`, launchd respawned it within 5s.

     Install per node (needs passwordless-or-already-authed sudo, which was
     confirmed available on the fleet): substitute __USER__, __VENV_PYTHON__
     (murakumo-venv's python3), __COMFYUI_DIR__ (~/comfyui), __HOME__, then:
       sudo cp com.murakumo.comfyui.plist /Library/LaunchDaemons/com.murakumo.comfyui.plist
       sudo chown root:wheel /Library/LaunchDaemons/com.murakumo.comfyui.plist
       sudo chmod 644 /Library/LaunchDaemons/com.murakumo.comfyui.plist
       sudo launchctl bootstrap system /Library/LaunchDaemons/com.murakumo.comfyui.plist
     Uninstall:
       sudo launchctl bootout system/com.murakumo.comfyui
       sudo rm /Library/LaunchDaemons/com.murakumo.comfyui.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>            <string>com.murakumo.comfyui</string>
  <key>UserName</key>         <string>__USER__</string>
  <key>ProgramArguments</key>
  <array>
    <string>__VENV_PYTHON__</string>
    <string>__COMFYUI_DIR__/main.py</string>
    <string>--listen</string>
    <string>0.0.0.0</string>
    <string>--port</string>
    <string>8188</string>
    <string>--disable-auto-launch</string>
  </array>
  <key>WorkingDirectory</key> <string>__COMFYUI_DIR__</string>
  <key>RunAtLoad</key>        <true/>
  <key>KeepAlive</key>        <true/>
  <key>ThrottleInterval</key> <integer>15</integer>
  <key>StandardOutPath</key>  <string>__HOME__/.gftd/comfyui.log</string>
  <key>StandardErrorPath</key><string>__HOME__/.gftd/comfyui.log</string>
</dict>
</plist>
