# deploy/murakumo-infer-compat.socket
[Unit]
Description=Murakumo inference tunnel compatibility listener

[Socket]
ListenStream=127.0.0.1:8090
NoDelay=true

[Install]
WantedBy=sockets.target
